# bgn_ngx/build/logagent/debian/debian/logagent.postinst
#!/bin/bash
# postinst script for logagent
#
# see: dh_installdeb(1)

###############################################################################
#
#   Copyright (C) Chaoyong Zhou
#   Email: bgnvendor@163.com
#   QQ: 2796796
#
################################################################################

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

docker_switch=off

prefix=/usr/local/logagent
service_dir=/etc/systemd/system

case "$1" in
    configure)
        echo "---------------------------------------------------------------------" >&2
        echo "  postinst $1  " >&2
        echo "---------------------------------------------------------------------" >&2

#        echo -n '[postinst] init logagent ...................... '
#        if [ "${docker_switch}" == "on" ]; then
#            chmod 0755 ${prefix}/bin/logagent_init.sh
#        else
#            if [ ! -f /etc/systemd/system/logagent.deb.service ]; then
#                echo "error: not found /etc/systemd/system/logagent.deb.service" >&2
#                exit 1
#            fi
#            #bash ${prefix}/bin/logagent_init.sh > /data/proclog/log/logagent/logagent_init.log 2>&1
#        fi
#
#        echo 'done'
#
        if [ -x "/etc/init.d/logagent" ]; then
            update-rc.d logagent defaults >/dev/null || exit $?
        fi

#        # reload due to service changed
#        systemctl daemon-reload
#
#        num=$(ls -1p ${service_dir}/logagent* > /dev/null)
#        if [ $? -eq 0 -a $num -gt 0 ]; then
#            for service_file in $(ls -1p ${service_dir}/logagent* > /dev/null)
#            do
#                service_file_des=$(basename ${service_file})
#    
#                systemctl enable ${service_file_des}
#    
#    #            if [ "${docker_switch}" != "on" ]; then
#    #                echo '[postinst] start ' ${service_file_des} ' ... '
#    #                systemctl start ${service_file_des}
#    #                echo '[postinst] start ' ${service_file_des} ' done '
#    #
#    #                #echo '[postinst] check ' ${service_file_des} ' status '
#    #                #systemctl status ${service_file_des}
#    #            fi
#            done
#        fi

        # ---- crontab ----
        echo "[postinst] restart cron"
        systemctl restart cron
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        echo "---------------------------------------------------------------------" >&2
        echo "  postinst $1  " >&2
        echo "---------------------------------------------------------------------" >&2

        if [ -x "/etc/init.d/logagent" ]; then
            update-rc.d logagent defaults >/dev/null || exit $?
        fi
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
